# moodle.service
[Unit]
Description=Moodle Bot Daemon

[Service]
Type=simple
WorkingDirectory=/home/alwin/bots/moodle
ExecStart=/home/alwin/bots/moodle/moodle_daemon.py
User=alwin

[Install]
WantedBy=multi-user.target
